# service/upgrade-aur.service
[Unit]
Description=Upgrades 'aur' repo
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/ctlos-aur/repo.sh upgrade
